# infra/systemd/pogo-lab-backup.service
[Unit]
Description=Backup diario de PostgreSQL para Pogo-lab
After=postgresql.service network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pogo-lab
EnvironmentFile=/opt/pogo-lab/.env
Environment=BACKUP_DIR=/var/backups/pogo-lab
Environment=BACKUP_RETENTION_DAYS=14
ExecStart=/opt/pogo-lab/bin/backup.sh

# Endurecimiento compatible con pg_dump y escritura en el directorio de backup.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/pogo-lab
